# deploy/com.cryptobot.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cryptobot.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ruiwang/Desktop/bian/.venv/bin/cryptobot</string>
        <string>daemon</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ruiwang/Desktop/bian</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ruiwang/Desktop/bian/logs/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ruiwang/Desktop/bian/logs/daemon.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ruiwang/Desktop/bian/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
